# components/cookbooks/postgresql-governor/templates/default/governor.service.erb
# Governor service unit file
#
# Auto generated by OneOps
# Created On : <%= Time.now.getutc %>

[Unit]
Description=Governor
Documentation=https://github.com/compose/governor
After=network.target
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/var/lib/pgsql
User=postgres
ExecStart=/bin/bash -c "PATH=/usr/pgsql-<%= node['postgresql-governor']['version'] %>/bin/:$PATH; export PATH; <%= node['postgresql-governor']['governor_home'] %>/governor.py <%= node['postgresql-governor']['governor_home'] %>/postgres.yml"
Restart=on-failure
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
