# components/cookbooks/postgresql-governor/templates/default/haproxy_status.service.erb
# haproxy_status service unit file
#
# Auto generated by OneOps
# Created On : <%= Time.now.getutc %>

[Unit]
Description=HAProxy_status
Documentation=https://github.com/compose/governor
After=network.target
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/var/lib/pgsql
User=postgres
ExecStart=/bin/bash -c "<%= node['postgresql-governor']['governor_home'] %>/haproxy_status.py <%= node['postgresql-governor']['governor_home'] %>/postgres.yml"
Restart=on-failure
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
