# components/cookbooks/etcd/templates/default/etcd.service.erb
# Etcd service unit file
#
# Auto generated by OneOps
# Created On : <%= Time.now.getutc %>

[Unit]
Description=Etcd Server
Documentation=https://coreos.com/etcd/docs
After=network.target
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
WorkingDirectory=<%= node.etcd.working_location%>
EnvironmentFile=<%= node.etcd.conf_file%>
User=root
# Set GOMAXPROCS to number of processors
ExecStart=/bin/bash -c "GOMAXPROCS=$(nproc) /usr/bin/etcd"
Restart=on-failure
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target