# components/cookbooks/docker_engine/templates/default/docker.service.erb
# Docker service unit file
#
# Auto generated by OneOps
# Created On : <%= Time.now.getutc %>

[Unit]
Description=Docker Application Container Engine
Documentation=https://docs.docker.com
After=network.target docker.socket
Requires=docker.socket

[Service]
Type=notify
ExecStart=/usr/bin/docker daemon -H fd://
MountFlags=slave
<% JSON.parse(node.docker_engine.limit_directives).each { |k, v| %>
<%= k %>=<%= v %>
<% } %>

[Install]
WantedBy=multi-user.target