# components/cookbooks/docker_engine/templates/default/docker.service.erb
# Docker service unit file
#
# Auto generated by OneOps
# Created On : <%= Time.now.getutc %>

[Unit]
Description=Docker Application Container Engine
Documentation=https://docs.docker.com
After=network.target docker.socket
Requires=docker.socket


[Service]
Type=notify
# The default is not to use systemd for cgroups because the delegate issues still
# exists and systemd currently does not support the cgroup feature set required
# for containers run by docker
ExecStart=/usr/bin/docker daemon -H fd://
MountFlags=slave
TimeoutStartSec=0
# set delegate yes so that systemd does not reset the cgroups of docker containers
Delegate=yes
<% JSON.parse(node.docker_engine.limit_directives).each { |k, v| %>
    <%= k %>=<%= v %>
<% } %>

[Install]
WantedBy=multi-user.target



