# components/cookbooks/graphite/files/default/uwsgi_graphite.conf
description "uwsgi"
start on runlevel [2345]
stop on runlevel [06]

exec /opt/uwsgi/sbin/uwsgi --log-syslog --die-on-term --emperor /opt/uwsgi/apps
